# server/deploy/lilyhub.service
[Unit]
Description=LilyHub local app server for LilyGo devices
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/lilyhub/lilyhub
Restart=on-failure
RestartSec=3
User=lilyhub
Group=lilyhub
Environment=LILYHUB_ADDR=:8080

# Writable state for uploaded notes. systemd creates /var/lib/lilyhub on
# start, chowns it to the service user, and implicitly adds it to
# ReadWritePaths — required because ProtectHome=yes masks $HOME (the
# default notes dir) and ProtectSystem=strict makes everything else
# read-only, so /api/notes/upload would otherwise 500 on mkdir/write.
StateDirectory=lilyhub
StateDirectoryMode=0750
Environment=LILYHUB_NOTES_DIR=/var/lib/lilyhub/notes

# Secrets (GROQ_API_KEY for /api/chat, etc.) live in /etc/lilyhub/env so they
# stay out of this unit file and out of git. The leading `-` makes the file
# optional — the service starts fine without it; only /api/chat returns 503.
# Recommended perms: chmod 0600, chown root:root. systemd reads it as PID 1
# before dropping privs to the lilyhub user, so root-only is correct.
#
# Setup (one-liner):
#   echo 'GROQ_API_KEY=gsk_xxx' | sudo tee /etc/lilyhub/env >/dev/null \
#     && sudo chmod 0600 /etc/lilyhub/env \
#     && sudo systemctl restart lilyhub
EnvironmentFile=-/etc/lilyhub/env

# Hardening — the binary only needs to bind a TCP port and make outbound HTTP.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
